# examples/interest-bot-upstart.conf
# interest-bot-upstart.conf
description "Crystal Interest Bot Upstart job"

start on runlevel [2345]
stop  on runlevel [!2345]

script
	cd /home/crystal-interest/
	bin/interest-bot-daemon
end script
